# _hosts/media.codonaft/etc/init.d/broadcastr
#!/sbin/openrc-run

command="/usr/local/bin/${RC_SVCNAME}"
command_user="nobody:nobody"
command_background="true"
command_args="${broadcastr_args}"
retry="10"

pidfile="/run/${RC_SVCNAME}.pid"
output_logger="logger -t ${RC_SVCNAME}"
error_logger="logger -t ${RC_SVCNAME}"

rc_ulimit="-n 65535"

depend() {
	need net
	use logger
}
