# _hosts/media.codonaft/etc/init.d/metasearch
#!/sbin/openrc-run

command="/usr/local/bin/${RC_SVCNAME}"
command_user="nobody:nobody"
command_background="true"
command_args="/etc/${RC_SVCNAME}/${RC_SVCNAME}.toml"

pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	need net
	use logger
}
