# _hosts/media.codonaft/etc/init.d/rnostr
#!/sbin/openrc-run

name="rnostr"

command="/usr/local/bin/${name}"
command_user="${name}:${name}"
command_background="true"
command_args="relay --watch -c /etc/${name}/${name}.toml"

pidfile="/run/${name}.pid"
output_logger="logger"
error_logger="logger"

depend() {
	need net
	use logger
}
