# _hosts/media.codonaft/etc/init.d/rnostr
#!/sbin/openrc-run

name="rnostr"

command="/usr/local/bin/${name}"
command_user="${name}:${name}"
command_background="true"
command_args="relay -c /etc/${name}/${name}.toml"
retry="60"

pidfile="/run/${name}.pid"
output_logger="logger -t ${name}"
error_logger="logger -t ${name}"

depend() {
	need net
	use logger
}
